# debian/git-daemon-sysvinit.postrm
#!/bin/sh
set -e

#DEBHELPER#

test "$1" = 'remove' || exit 0

update-rc.d git-daemon remove >/dev/null
